# iis.service
[Unit]
Description=Intelligence Integration System Services
After=network.target

[Service]
Type=simple
User=sleepy
Group=sleepy
WorkingDirectory=/home/sleepy/Documents/IntelligenceIntegrationSystem
ExecStart=/home/sleepy/Documents/IntelligenceIntegrationSystem/launch_iis_system.sh
Restart=always
RestartSec=10

# 保证日志也能从 journalctl 看到
StandardOutput=journal
StandardError=journal

# 停止时先发 SIGTERM，超时后再强制停止
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

